# contrib/init/clashicd.conf
description "Bitcoin Clashic ABCD Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env CLASHICD_BIN="/usr/bin/clashicd"
env CLASHICD_USER="bitcoin"
env CLASHICD_GROUP="bitcoin"
env CLASHICD_PIDDIR="/var/run/clashicd"
# upstart can't handle variables constructed with other variables
env CLASHICD_PIDFILE="/var/run/clashicd/clashicd.pid"
env CLASHICD_CONFIGFILE="/etc/clashic/clashic.conf"
env CLASHICD_DATADIR="/var/lib/clashicd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # clashicd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$CLASHICD_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run clashicd."
        echo "The setting must appear in $CLASHICD_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Bitcoin Clashic Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$CLASHICD_PIDDIR"
    chmod 0755 "$CLASHICD_PIDDIR"
    chown $CLASHICD_USER:$CLASHICD_GROUP "$CLASHICD_PIDDIR"
    chown $CLASHICD_USER:$CLASHICD_GROUP "$CLASHICD_CONFIGFILE"
    chmod 0660 "$CLASHICD_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$CLASHICD_PIDFILE" \
    --chuid $CLASHICD_USER:$CLASHICD_GROUP \
    --exec "$CLASHICD_BIN" \
    -- \
    -pid="$CLASHICD_PIDFILE" \
    -conf="$CLASHICD_CONFIGFILE" \
    -datadir="$CLASHICD_DATADIR" \
    -disablewallet \
    -daemon

